# apploader/conf/appload-supervisord.ini
[program: apploader]
command=/workplace/apploader/apploader     
directory=/workplace/apploader            ;c
priority=1                    ;9991
numprocs=1                    ;
autostart=true                ;supervisor
autorestart=unexpected        ; exitcodes
startretries=10               ;
exitcodes=0                   ;）
stopsignal=KILL               ;
stopwaitsecs=10               ;SIGKIL
redirect_stderr=true          ;stderstdout
stopasgroup=true              ; kil
